# scheduler/cups.service.in
[Unit]
Description=CUPS Scheduler
Documentation=man:cupsd(8)
After=network.target nss-user-lookup.target time-set.target nslcd.service
Requires=cups.socket

[Service]
ExecStart=@sbindir@/cupsd -l
Type=notify
Restart=on-failure
Slice=system-cups.slice
@TIMEOUTSTARTSEC@

[Install]
Also=cups.socket cups.path
WantedBy=@SYSTEMD_WANTED_BY@
