# scheduler/cups.socket.in
[Unit]
Description=CUPS Scheduler
PartOf=cups.service

[Socket]
ListenStream=@CUPS_DEFAULT_DOMAINSOCKET@
RemoveOnStop=on

[Install]
WantedBy=sockets.target
