# scheduler/system-cups.slice
[Unit]
Description=CUPS Slice
Documentation=man:cups(1)
